# templates/default/epipe.service.erb
[Unit]
Description = ePipe Server
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=120
StartLimitBurst=15

[Service]
User = <%= node['epipe']['user'] %>

PIDFile = <%= node['epipe']['pid_file'] %>
ExecStart = <%= node['epipe']['base_dir'] %>/bin/start-epipe.sh
ExecStop = <%= node['epipe']['base_dir'] %>/bin/stop-epipe.sh
LimitNOFILE = 20000
Restart=on-failure
RestartSec=5

Type = forking

[Install]
WantedBy = multi-user.target
